# resources/info.timer
[Unit]
Description=AKTIN Update-Info Service Timer

[Timer]
OnCalendar=hourly

[Install]
WantedBy=basic.target

