# resources/update.service
[Unit]
Description=AKTIN Update Service
After=network.target

[Service]
ExecStart=__PACKAGE__

[Install]
WantedBy=default.target

