# debian/prerm
#!/bin/bash
set -euo pipefail

systemctl stop __PACKAGE__.socket __PACKAGE__.service __PACKAGE__-info.timer __PACKAGE__-info.service
systemctl disable __PACKAGE__.socket __PACKAGE__-info.timer

